# old_exams_done/2021_09_06/vagrant_ansible/roles/add_systemd_unit/files/unit.service
[Service]
Type=simple
Description=Run count.sh in background to log ssh connections coming from clients
ExecStart=/home/count.sh
KillMode=process
Restart=on-failure

[Install]
WantedBy=basic.target
